# deploy/systemd/lite-api-server.service
[Unit]
Description=Lite API Server
After=network.target

[Service]
Type=simple
User=lite-api
WorkingDirectory=/opt/lite-api-server
ExecStart=/opt/lite-api-server/lite-api-server --config /etc/lite-api-server/config.yaml
Restart=always
RestartSec=3
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
